# upstart/sakamotosan.conf
description "sakamotosan"
author  "yoshikyoto"

start on runlevel [2345]
stop on runlevel [016]

chdir /home/sakamotosan/sakamotosan
exec yarn run start
respawn
